# configs/antergos/root-image/etc/systemd/system/pacman-boot.service
[Unit]
Description=Sync pacman on first boot
After=pacman-init.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/pacman-boot