# configs/cinnarch_minimum/root-image/usr/lib/systemd/system/livecd.service
[Unit]
Description=LiveMedia Config Script
Before=systemd-user-sessions.service

[Service]
Type=oneshot
ExecStart=/etc/cinnarch/livecd

[Install]
WantedBy=multi-user.target
