# configs/openbox/back/root-image/etc/systemd/system/getty@tty1.service.d/autologin.conf
[Service]
ExecStart=
ExecStart=-/usr/bin/agetty --autologin antergos --noclear %I 38400 linux
